# lablivre-4.1.8/lablivre-4.1.8/packaging/debian/DEBIAN/postrm
#!/bin/bash
set -e

# ============================================================
# postrm — Após remover o pacote
# ============================================================

# Sempre remove os symlinks (em remove e purge)
rm -f /usr/local/bin/lablivre
rm -f /usr/local/bin/lablivre-gui

# Descobre o usuário operador (mesmo critério do postinst)
USUARIO_REAL="${SUDO_USER:-}"
if [ -z "$USUARIO_REAL" ] || [ "$USUARIO_REAL" = "root" ]; then
    USUARIO_REAL=$(awk -F: '$3 >= 1000 && $3 < 65534 {print $1; exit}' /etc/passwd)
fi

# Remove serviço systemd do servidor web (em qualquer remoção)
systemctl disable lablivre-web.service 2>/dev/null || true
systemctl stop lablivre-web.service 2>/dev/null || true
rm -f /etc/systemd/system/lablivre-web.service
rm -f /etc/systemd/system/wol-*.service
systemctl daemon-reload 2>/dev/null || true

# Remove as tarefas do cron do usuário operador
if [ -n "$USUARIO_REAL" ] && id "$USUARIO_REAL" >/dev/null 2>&1; then
    crontab -u "$USUARIO_REAL" -l 2>/dev/null | grep -v "# LabLivre" | \
        crontab -u "$USUARIO_REAL" - 2>/dev/null || true
fi

# Em PURGE (remoção completa): apaga tudo, inclusive dados gerados e configs
if [ "$1" = "purge" ]; then
    rm -rf /opt/lablivre
    echo "[postrm] LabLivre removido completamente (purge)."
else
    # Em remoção simples: mantém configs e dados, mas avisa
    echo "[postrm] LabLivre removido."
    echo "          Configurações e logs preservados em /opt/lablivre/"
    echo "          Para apagar tudo: sudo apt purge lablivre && sudo rm -rf /opt/lablivre"
fi

exit 0
